# linux/geometry-ai-server_1.0.0.90_amd64/etc/systemd/system/geometry-ai.service
[Unit]
Description=Geometry AI Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/geometry-ai
ExecStart=/usr/bin/python3 /usr/local/geometry-ai/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
