# geometry-ai-server_1.0.0.0707_amd64/DEBIAN/prerm
#!/bin/bash
systemctl stop geometry-ai 2>/dev/null || true
systemctl disable geometry-ai 2>/dev/null || true
